# backend/app/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi.ini
stderr_logfile=uwsgi_err.log
stdout_logfile=uwsgi_out.log

[program:nginx-app]
command = /usr/sbin/nginx
